# board/opendingux/target_skeleton/etc/profile.d/banner.sh
if [ "$-" ]; then
	cat /etc/welcome-banner.txt
fi
