# board/opendingux/target_skeleton/etc/profile.d/love2d.sh
# Force GLES graphics for games created with Love2D
export LOVE_GRAPHICS_USE_OPENGLES=1
